# infra/mdns/com.clawconsole.mdns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon that advertises claw.local over mDNS at boot and keeps it alive.

  Install (root, runs before login so it survives reboots):
      sudo cp infra/mdns/com.clawconsole.mdns.plist /Library/LaunchDaemons/
      sudo launchctl load -w /Library/LaunchDaemons/com.clawconsole.mdns.plist

  The ProgramArguments path is absolute — if this repo moves, update it and
  reload. Logs go to /tmp/claw-mdns.log.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawconsole.mdns</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fernandolincoln/repos/claw-console/infra/mdns/advertise-claw-local.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claw-mdns.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claw-mdns.log</string>
</dict>
</plist>
